# templates/default/debian-mongodb.upstart.erb
#
# Automatically Generated by Chef, do not edit directly!
#

<%- @ulimit.to_hash.keys.sort.each do |key| %>
    <%- if not @ulimit[key].nil? %>
limit <%= key %> <%= @ulimit[key] %> <%= @ulimit[key] %>
    <% end %>
<%- end %>

kill timeout 300 # wait 300s between SIGTERM and SIGKILL.

start on runlevel [2345]
stop on runlevel [06]

script
  NAME=<%= @provides %>
  ENABLE_MONGOD="yes"
  if [ -f <%= @sysconfig_file %> ]; then
    . <%= @sysconfig_file %>;
  fi
  if [ "x$ENABLE_MONGOD" = "xyes" ]; then
      #Mongo sometimes losts his pid
      dbpath=`cat $CONFIGFILE | egrep ^dbpath | awk -F '=' '{print $2}'`;
      mkdir -p $dbpath
      chown $DAEMON_USER: $dbpath
      if [ -e $dbpath/mongod.lock ] ; then
          pid=`cat $dbpath/mongod.lock`;
          cmd=`cat /proc/$pid/cmdline | tr "\000" "\n"|head -n 1 |cut -d : -f 1`
          if [ "$cmd" == "/usr/bin/mongod" ]; then
              echo "Mongodb already start";
              exit 1;
          else
              echo "Hmm..old pid.. i remove it";
              rm  $dbpath/mongod.lock;
          fi;
      fi
      pidfilepath=`cat $CONFIGFILE | egrep ^pidfilepath | awk -F '=' '{print $2}'`;
      pidfiledir=`basename $pidfilepath`
      mkdir -p $pidfiledir
      chown $DAEMON_USER: $pidfiledir
      exec start-stop-daemon --start --pidfile $pidfilepath --quiet --chuid $DAEMON_USER --exec $DAEMON -- $DAEMON_OPTS
  fi
end script
